# bot/avatar-queue-bot.service
# Установка на VM (как у Profile Mirror):
#   sudo cp avatar-queue-bot.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now avatar-queue-bot
#   journalctl -u avatar-queue-bot -f   # логи

[Unit]
Description=Telegram bot: avatar video queue
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=xen
WorkingDirectory=/opt/avatar-pipeline/bot
# Секреты храним в отдельном файле с правами 600:
#   /opt/avatar-pipeline/bot/.env  (KEY=value построчно)
EnvironmentFile=/opt/avatar-pipeline/bot/.env
ExecStart=/usr/bin/python3 /opt/avatar-pipeline/bot/bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
